# ansible/roles/services_data/templates/garage.service

[Unit]
Description=Garage S3-compatible object store
After=network-online.target
Wants=network-online.target

[Service]
User=garage
Group=garage
ExecStart=/usr/local/bin/garage --rpc-secret-file /etc/garage/rpc_secret --admin-token-file /etc/garage/admin_token --metrics-token-file /etc/garage/metrics_token -c /etc/garage.toml server --single-node
Restart=on-failure
RestartSec=5s
LimitNOFILE=65536
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/garage

[Install]
WantedBy=multi-user.target
